# systemd/system/upload-pictures.path
[Unit]
Description=Watch the Pictures directory for changes
StartLimitIntervalSec=3min

[Path]
PathChanged=/home/pi-camera/Pictures
PathExistsGlob=/home/pi-camera/Pictures/*.jpg
MakeDirectory=true
#TriggerLimitIntervalSec=3min

[Install]
WantedBy=multi-user.target
