# systemd/user/rclone-upload-pictures.path
[Unit]
Description=Watch the Pictures directory for changes

[Path]
PathChanged=%h/Pictures
MakeDirectory=true
TriggerLimitIntervalSec=3min

[Install]
WantedBy=default.target
