# systemd/user/upload-pictures.path
[Unit]
Description=Watch the Pictures directory for changes
StartLimitIntervalSec=3min

[Path]
PathChanged==%h/Pictures
PathExistsGlob=%h/Pictures/*.jpg
MakeDirectory=true
#TriggerLimitIntervalSec=3min

[Install]
WantedBy=default.target
